# deploy/lastsys.service
[Unit]
Description=LastSys Spring Boot application
After=network-online.target
Wants=network-online.target

[Service]
User=lastsys
Group=lastsys
WorkingDirectory=/opt/lastsys
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=SPRING_DATASOURCE_URL=jdbc:mysql://127.0.0.1:3306/lastsys?serverTimezone=Asia/Seoul&characterEncoding=UTF-8
Environment=SPRING_DATASOURCE_USERNAME=lastsys
Environment=SPRING_DATASOURCE_PASSWORD=change-me
ExecStart=/usr/bin/java -jar /opt/lastsys/app.jar
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
